# files/systemd/nix.mount
# Bind-mount the persistent, per-machine /var/nix onto /nix.
#
# On this immutable-root bootc system /nix is not writable, and Nix refuses a
# symlinked /nix. The binary cache (cache.nixos.org) ships artifacts hard-coded
# to absolute /nix/store paths, so the store must live at a real /nix for cache
# substitution to work. /var is machine-local, writable, and persists across
# image upgrades, so we keep the store in /var/nix and bind it onto /nix.
#
# The unit filename MUST be nix.mount to match the mountpoint /nix.
#
# nix-store-init.service seeds /var/nix and must run first. Requiring and
# ordering after it fixes the first-boot race that the old
# ConditionPathExists=/var/nix guard caused: /var/nix did not yet exist when the
# mount was evaluated, so it was silently skipped and never retried. See
# docs/nix-store-boot-race.md.
[Unit]
Description=Bind mount /var/nix onto /nix
Before=nix-daemon.service nix-daemon.socket
Requires=nix-store-init.service
After=nix-store-init.service

[Mount]
What=/var/nix
Where=/nix
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
