# files/systemd/nix-daemon.service.d/10-nix-mount.conf
# Ensure the Nix store is mounted before the daemon starts. /var/nix is created
# on first boot by tmpfiles and bind-mounted onto /nix by nix.mount.
[Unit]
RequiresMountsFor=/nix
Requires=nix.mount
After=nix.mount
